# old/dz3_4_OS_part2/node_exporter.service
[Unit]
Description=Metrics collector

[Service]
EnvironmentFile=-./nodeexporter.conf
ExecStart=/home/vagrant/node_exporter/node_exporter -f $EXTRA_OPTS
[Install]
WantedBy=multi-user.target